# docs/archive/linux/slt-dashboard.service
[Unit]
Description=SLT Call Center — Transcription Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=slt
WorkingDirectory=/opt/slt/Telecom-Voice-to-Text
ExecStart=/opt/slt/Telecom-Voice-to-Text/.venv/bin/python dashboard_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
